# deploy/systemd/2dehands-monitor.service
[Unit]
Description=2dehands.be new listing monitor (single check cycle)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=2dehands-monitor
Group=2dehands-monitor
WorkingDirectory=/opt/2dehands-monitor
EnvironmentFile=/opt/2dehands-monitor/.env
ExecStart=/opt/2dehands-monitor/.venv/bin/python /opt/2dehands-monitor/monitor.py

# Light sandboxing - this process only needs outbound HTTPS and its own
# data/log files, nothing else on the system.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/2dehands-monitor/data
